# chef-repo/cookbooks/screeningtest/templates/upstart.conf.erb
# screeningtest - golang http server

description	"golang http server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /opt/screeningtest/screeningtest
